# share/crontab
@hourly ${HOME}/.local/bin/process-review-requested-notifications
