# share/crontab
SHELL=/usr/local/bin/fish

@hourly process-review-requested-notifications
