# share/crontab
0 9,13 * * 1-5 ${HOME}/.local/bin/process-review-requested-notifications
