# share/crontab
0 9,14 * * 1-4 ${HOME}/.local/bin/process-review-requested-notifications
